# tools/princess_nightly/crontab
0 5 * * * /usr/local/bin/nightly.sh dev 
15  5 * * * /usr/local/bin/nightly.sh master
